# systemd/ai-chat-server.service
[Unit]
Description=AI MCP Chat Server for operation-dbus
Documentation=https://github.com/operation-dbus
After=network.target

[Service]
Type=simple
User=jeremy
WorkingDirectory=/git/operation-dbus
Environment="NODE_ENV=production"
Environment="PATH=/home/jeremy/.nvm/versions/node/v22.19.0/bin:/usr/local/bin:/usr/bin:/bin"
# Set default model via environment variable
Environment="OLLAMA_DEFAULT_MODEL=llama2"
ExecStart=/home/jeremy/.nvm/versions/node/v22.19.0/bin/node /git/operation-dbus/chat-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-chat

[Install]
WantedBy=multi-user.target
